# rc.local
# Add following line to rc.local

sudo python /home/pi/projects/RPi_CPU_Fan/cpufan.py &
