# rc.local
# Add following line to rc.local

sudo python /home/pi/projects/fan/cpufan.py &
